# provisioning/roles/key-sync/files/pull-keys@.timer
[Unit]
Description=Update authorized SSH keys for %i

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min

[Install]
WantedBy=timers.target
