# provisioning/roles/unbound/templates/unbound-update-block-list.timer
[Unit]
Description=Update unbound block lists on a schedule

[Timer]
OnCalendar=daily
RandomizedDelaySec=3600

[Install]
WantedBy=timers.target
