# provisioning/roles/key-sync/templates/pull-keys@.timer
[Unit]
Description=Update authorized SSH keys for %i

[Timer]
OnUnitActiveSec=15m

[Install]
WantedBy=timers.target
